# systemd/rfserver.service
[Unit]
Description=Create Signed Urls Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
ExecStart=/usr/local/bin/rfserve
Restart=always
RestartSec=5
EnvironmentFile=/etc/rfserve/.env
WorkingDirectory=/var/lib/rfserve

[Install]
WantedBy=multi-user.target